# cexio.conf
[program:cexio]
command=docker-compose -f production.yml up
directory=/cexio
redirect_stderr=true
autostart=true
autorestart=true
priority=10